# zookeeper/zk_supervisor.conf
[supervisord]
nodaemon = true

[program:zookeeper-1]
command = bin/zookeeper-server-start.sh config/zookeeper.properties
directory = /kafka