# broker/broker_supervisor.conf
[supervisord]
nodaemon = true

[program:broker]
command = bin/kafka-server-start.sh config/server.properties
directory = /kafka